# rel/service/linux/dust.service
[Unit]
Description=Dust Distributed File System Daemon
Documentation=https://github.com/AndrewBoessen/dust
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/usr/local/bin/dust start
ExecStop=/usr/local/bin/dust stop
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/dust
PrivateTmp=yes

# Resource limits
LimitNOFILE=65536

# Environment
Environment=HOME=/var/lib/dust
Environment=DUST_DATA_DIR=/var/lib/dust
Environment=RELEASE_NODE=dust@127.0.0.1

[Install]
WantedBy=multi-user.target
