# apps/dust_api/priv/service/linux/dust.service
[Unit]
Description=Dust Distributed File System Daemon
Documentation=https://github.com/AndrewBoessen/dust
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/opt/dust/bin/dust start
ExecStop=/opt/dust/bin/dust stop
Restart=on-failure
RestartSec=5

# Resource limits
LimitNOFILE=65536

# Environment
Environment=RELEASE_NODE=dust@127.0.0.1

[Install]
WantedBy=multi-user.target
